# rating-website.service
[Unit]
Description=Skjærvik's Ratings Website
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/your/project
ExecStart=/usr/bin/npm start
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
